# debian/postinst
#!/bin/bash
set -e

# Cria user de sistema se não existir
if ! id -u ishin-gateway >/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin ishin-gateway
fi

# Cria diretórios de runtime
mkdir -p /var/log/ishin-gateway
mkdir -p /var/lib/ishin-gateway
mkdir -p /etc/ishin-gateway/ssl
mkdir -p /etc/ishin-gateway/rules

# Ajusta permissões
chown -R ishin-gateway:ishin-gateway /var/log/ishin-gateway
chown -R ishin-gateway:ishin-gateway /var/lib/ishin-gateway
chown -R ishin-gateway:ishin-gateway /etc/ishin-gateway/ssl
chown -R ishin-gateway:ishin-gateway /etc/ishin-gateway/rules

# Habilita o serviço (não inicia automaticamente)
systemctl daemon-reload
systemctl enable ishin-gateway.service

echo "ishin-gateway installed successfully."
echo "  - Edit /etc/ishin-gateway/adapter.yaml to configure"
echo "  - Place rules in /etc/ishin-gateway/rules/"
echo "  - Start with: systemctl start ishin-gateway"
